# config/saba-chan.service
# config/saba-chan.service
[Unit]
Description=Saba-chan Game Server Manager Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=saba-chan
Group=saba-chan
WorkingDirectory=/opt/saba-chan
ExecStart=/opt/saba-chan/saba-core
Restart=on-failure
RestartSec=10
Environment="SABA_DATA_DIR=/opt/saba-chan/data"

# 보안 강화
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/saba-chan

[Install]
WantedBy=multi-user.target
